# assets/s6/rtorrent/run
#!/bin/sh
echo "Spawning rTorrent"
rtorrent -I -D -n -o import=/home/rtorrent/.rtorrent.rc > /var/log/s6/rtorrent 2>&1