# assets/s6/rtorrent/run
#!/bin/sh
echo "Spawning rTorrent"

rtorrent -I -D -n -o import=/rtorrent/.rtorrent.rc > rtorrentstart 2>&1